# wikiuzru.conf
[program:wikiuzbot]
directory=/home/ubuntu/wikipediauzbot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/wikipediauzbot/logfile_err.log
stdout_logfile=/home/ubuntu/wikipediauzbot/logfile.log